# init.d/bootmisc
#!/sbin/runscript
# Copyright (c) 2007-2009 Roy Marples <roy@marples.name>
# Copyright © 2007-2018 ANSSI. All Rights Reserved.
# Modified for CLIP: Vincent Strubel <clipos@ssi.gouv.fr>
# Released under the 2-clause BSD license.

depend()
{
	need localmount
	before logger
	after clock sysctl
	keyword -timeout
}

: ${log_dmesg:=${LOG_DMESG:-yes}}

mkutmp()
{
	: >"$1"
	# Not all systems have the utmp group
	chgrp utmp "$1" 2>/dev/null
	chmod 0664 "$1"
}

start()
{
	ebegin "Creating user login records"
	local xtra=
	for x in "" $xtra; do
		mkutmp /var/run/utmp$x
	done
	[[ -e /var/log/wtmp ]] || mkutmp /var/log/wtmp
	eend 0

	ebegin "Cleaning /var/run and /var/lock"
	for x in $(find /var/run ! -type d ! -name utmp \
		! -name random-seed ! -name logrotate.conf);
	do
		rm -f -- "$x"
	done
	find /var/lock -type f -print0 | xargs -0 rm -f --
	eend 0

	# Make sure our X11 stuff have the correct permissions
	# Omit the chown as bootmisc is run before network is up
	# and users may be using lame LDAP auth #139411
	mkdir -p /tmp/.ICE-unix /tmp/.X11-unix
	chmod 1777 /tmp/.ICE-unix /tmp/.X11-unix
	if [[ -x /sbin/restorecon ]]; then
		restorecon /tmp/.ICE-unix /tmp/.X11-unix
	fi

	return 0
}

stop()
{
	# Write a halt record if we're shutting down
	if [[ "$RC_RUNLEVEL" == shutdown ]]; then
		[[ "$RC_UNAME" == Linux ]] && halt -w
	fi

	return 0
}

# vim: ft=sh
